# custom-pi-boot 2/config/etc/rc.local
#!/bin/bash
fbi -T 1 -d /dev/fb0 -noverbose -a /boot/splash.png &
cd /home/pi
./ebike_demo &
exit 0
